# scripts/systemd/news-crawler-ai.service
[Unit]
Description=News Crawler AI pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=christianfita
WorkingDirectory=/home/christianfita/apps/news-crawler-ai
Environment=CHROME_BIN=/snap/bin/chromium
Environment=CHROMEDRIVER_PATH=/usr/local/bin/chromedriver-143.0.7499.40
EnvironmentFile=/home/christianfita/apps/news-crawler-ai/.env
ExecStart=/home/christianfita/apps/news-crawler-ai/.venv/bin/python -m ingest.classifier
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
